# m0/backup/sas-m0-backup-placeholder.service
[Unit]
Description=Super Ainux backup placeholder (no-op; replace in later milestone)

[Service]
Type=oneshot
ExecStart=/bin/true
